# misc/drcomd.conf

description "drcomd service"

start on started lightdm
stop on shutdown


exec /usr/local/bin/drcomd

post-stop script
    modprobe -r drcom
end script
