# misc/drcomc.conf

description "drcomc login service"

start on started drcomd 
stop on shutdown


exec /usr/local/bin/drcomc login
